# rpi-led-demo/systemd/led-watcher.service
[Unit]
Description=GPIO LED Watcher

[Service]
Type=simple
ExecStart=/usr/bin/sh /opt/rpi-led-watcher/start.sh

[Install]
WantedBy=multi-user.target